# deploy/systemd/stylino-docker-compose.service
[Unit]
Description=Stylino (Docker Compose)
Requires=docker.service
After=docker.service network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=/srv/stylino
ExecStart=/usr/bin/docker compose --project-name stylino -f deploy/docker-compose.prod.yml up -d
ExecStop=/usr/bin/docker compose --project-name stylino -f deploy/docker-compose.prod.yml down
TimeoutStartSec=0

[Install]
WantedBy=multi-user.target
